# supervisord.conf
[supervisord]
nodaemon=true

[program:prov7adm]
command=service prov7adm start

[program:prov7dep]
command=service prov7dep start